# crontab
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/10 * * * * cd /grades && python3 pushNewGrades.py >> /grades/cron.log &2>1
